# supervisord.scheduler.conf
[supervisord]
nodaemon=true
user=root

[program:scheduler]
command=supercronic /etc/supercronic/supercronic
user=app
autostart=true
autorestart=true